# services/rabbitmq.service.conf
[program:rabbitmq]
command=RABBITMQ_CONFIG_FILE=/storage/rabbitmq.config /usr/sbin/rabbitmq-server
stdout_events_enabled = true
stderr_events_enabled = true
